# dot_config/containers/systemd/user/nextcloud.container
[Unit]
Description=Nextcloud Personal Cloud Server
After=network.target

[Container]
Image=nextcloud:latest
ContainerName=nextcloud
Volume=/run/media/tom/LaCie/nextcloud:/var/www/html:Z
AddDevice=/run/media/tom/LaCie
Environment=TZ=Etc/UTC
Environment=NEXTCLOUD_ADMIN_USER=admin
Environment=NEXTCLOUD_ADMIN_PASSWORD=password
PublishPort=8443:443

[Service]
TimeoutStartSec=900
Restart=on-failure

[Install]
WantedBy=default.target
