# dot_config/containers/systemd/user/nextcloud.container
[Unit]
Description=Nextcloud Personal Cloud Server
After=network.target

[Container]
Image=lscr.io/linuxserver/nextcloud:latest
ContainerName=nextcloud
Volume=/run/media/tom/LaCie/nextcloud/config:/config:Z
Volume=/run/media/tom/LaCie/nextcloud/data:/data:Z
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Etc/UTC
PublishPort=8443:443

[Service]
TimeoutStartSec=900
Restart=on-failure

[Install]
WantedBy=default.target
